# config/docker/supervisord_dev.conf
[supervisord]
nodaemon=true
logfile=/app/log/supervisord.log
childlogdir=/app/log
pidfile=/app/tmp/supervisor.pid

[program: xvfm]
command=

[program: rackup]
command=bundle exec rackup
